# container-files/etc/supervisor.d/httpd.conf
[program:httpd]
command=httpd -D "FOREGROUND" -k start
redirect_stderr=true
